# money-flow-momentum-live.service
# Live guarded instance: operating_mode=live_guarded, BYBIT_ENV=live
# Uses config/config.live.yaml, .env.live, data/live/, artifacts/live/, logs/live/
# Run simultaneously with money-flow-momentum-demo.service on the same host.
[Unit]
Description=Money Flow Momentum - Live guarded instance
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=REPO_USER
Group=REPO_GROUP
WorkingDirectory=REPO_ROOT
ExecStart=REPO_ROOT/venv/bin/python run_bot.py run --config config/config.live.yaml
Restart=on-failure
RestartSec=10
Environment=PATH=REPO_ROOT/venv/bin
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-REPO_ROOT/.env.live

StandardOutput=append:REPO_ROOT/logs/live/bot.log
StandardError=append:REPO_ROOT/logs/live/bot.log

[Install]
WantedBy=multi-user.target
