# money-flow-momentum-automation.timer
[Unit]
Description=Money Flow Momentum - Run Demo orchestration periodically
Documentation=Run "python run_bot.py automation cycle" every 15 minutes when enabled.

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true
Unit=money-flow-momentum-automation.service

[Install]
WantedBy=timers.target
